# services/systemd/pulse-health.service
[Unit]
Description=Pulse Health Monitor
After=network.target pulse-hub.service
Requires=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/pulse
Environment="PYTHONPATH=/opt/pulse"
ExecStart=/opt/pulse/venv/bin/python3 -m services.sensors.health_monitor
Restart=always
RestartSec=60
TimeoutStartSec=30

# Logging
StandardOutput=append:/var/log/pulse/health.log
StandardError=append:/var/log/pulse/health-error.log

[Install]
WantedBy=multi-user.target
